# supervisord.conf
[supervisord]
nodaemon=true

[program:nodeapp]
command=node server.js
directory=/app
