# supervisord.conf
[supervisord]
nodaemon=true

[program:serf]
command=serf agent -config-dir=/opt/config/startup -snapshot=/tmp/serf.snapshot

[program:nodejs_app]
command=nodejs server.js
directory=/app

