# supervisord.conf
[supervisord]
nodaemon=true

[program:this_is_some_junk]
command=node server.js
directory=/app
