# supervisord.conf
[supervisord]
nodaemon=true

[program:app1]
command=node server.js
directory=/app
